# distributed-dcompose-stack/data/lib/systemd/radar-output.timer.template
[Unit]
Description=RADAR-Docker hourly data restructuring
BindsTo=radar-docker.service

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h
Unit=radar-output.service

[Install]
WantedBy=default.target
