# dcompose-stack/radar-cp-hadoop-stack/hash-backup/systemd/radar-hashbackup.service.template
[Unit]
Description=RADAR-Docker hashbackup service

[Service]

WorkingDirectory=

ExecStart=

NotifyAccess=all
