# distributed-dcompose-stack/commons/lib/systemd/radar-check-health.service.template
[Unit]
Description=RADAR-Docker Health Check service

[Service]

WorkingDirectory=

ExecStart=

NotifyAccess=all

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=radar-check-health
