# distributed-dcompose-stack/commons/lib/systemd/radar-check-health.timer.template
[Unit]
Description=RADAR-Docker health checking every 5 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=5min
Unit=radar-check-health.service

[Install]
WantedBy=default.target