# dcompose-stack/radar-cp-hadoop-stack/lib/systemd/radar-output.service.template
[Unit]
Description=RADAR-Docker service
After=radar-docker.service
Requires=docker.service

[Service]

WorkingDirectory=

ExecStart=

NotifyAccess=all

