# distributed-dcompose-stack/proxies/lib/systemd/radar-renew-certificate.service.template
[Unit]
Description=RADAR-Docker renew certificate service

[Service]

WorkingDirectory=

ExecStart=

NotifyAccess=all

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=radar-renew-certificate
